# .crontab
0 * * * * /bin/bash -li -c "maid clean --force --silent"
