# .crontab
0 * * * * /bin/bash -li -c "maid clean --force --silent"
*/30 * * * * /Users/tor/.bin/sync-pass
